# packaging/vice.service
[Unit]
Description=Vice game clip recorder daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/vice start --no-open-ui
Restart=on-failure
RestartSec=3
PassEnvironment=WAYLAND_DISPLAY DISPLAY XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS XDG_SESSION_TYPE XDG_CURRENT_DESKTOP

[Install]
WantedBy=graphical-session.target
